# src-tauri/macos/helper/ru.ravel.ultunnel-macos.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ru.ravel.ultunnel-macos.helper</string>

  <key>MachServices</key>
  <dict>
    <key>ru.ravel.ultunnel-macos.helper</key>
    <true/>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
